# systemd/zplug-update.timer

[Unit]
Description=Run zplug update weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=default.target
