# systemd/zsh-compdump.timer
[Unit]
Description=Rebuild zcompdump daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
